# container-files/etc/dnsmasq.conf
[program:dnsmasq]
command=/usr/sbin/dnsmasq --no-daemon
priority=9
autostart=true
autorestart=true
